# Alice/ip.service
[Unit]
Description = IP script

[Service]
ExecStart=/home/user/Documents/ip.sh

[Install]
WantedBy=multi-user.target